# examples/com.matthias.mailfrost.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <!-- Load with:
         launchctl unload ~/Library/LaunchAgents/com.matthias.mailfrost.plist 2>/dev/null
         launchctl load ~/Library/LaunchAgents/com.matthias.mailfrost.plist
         launchctl start com.matthias.mailfrost

         Check with:
         launchctl list | grep mailfrost
    -->

    <key>Label</key>
    <string>com.matthias.mailfrost</string>

    <!-- Use absolute paths for the binary, config, and .env. -->
    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/bin/mailfrost</string>
      <string>-config</string>
      <string>/Users/matthias/backup-mail-m/config</string>
      <string>-env</string>
      <string>/Users/matthias/backup-mail-m/.env</string>
      <string>backup</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/matthias/backup-mail-m</string>

    <!-- StartInterval is in seconds. 7200 means every 2 hours. -->
    <key>StartInterval</key>
    <integer>7200</integer>

    <key>RunAtLoad</key>
    <true/>

    <!-- Mailfrost already takes a lock in state_dir/.lock, so overlapping scheduled runs fail safely.
         Review data/state/logs/ and alerts.log as part of normal monitoring.
    -->
    <key>StandardOutPath</key>
    <string>/Users/matthias/backup-mail-m/data/state/logs/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/matthias/backup-mail-m/data/state/logs/launchd.err.log</string>
  </dict>
</plist>
