# examples/com.matthias.mailfrost.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.matthias.mailfrost</string>

    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/bin/mailfrost</string>
      <string>-config</string>
      <string>/Users/matthias/backup-mail-m/config</string>
      <string>-env</string>
      <string>/Users/matthias/backup-mail-m/.env</string>
      <string>backup</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/matthias/backup-mail-m</string>

    <key>StartInterval</key>
    <integer>7200</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/matthias/backup-mail-m/data/state/logs/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/matthias/backup-mail-m/data/state/logs/launchd.err.log</string>
  </dict>
</plist>
